# deploy/cncflow-parser.service
[Unit]
Description=cncflow engineering file parser
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/dev/cncflow/backend
Environment=CNCFLOW_FILE_STORAGE=/var/lib/cncflow/files
ExecStart=/root/dev/cncflow/backend/.venv-parser/bin/python -m cncflow_core.ingestion.worker
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
